# deploy/crontab
# 海洋気象データ自動取得 - 10分間隔
*/10 * * * * cd /var/www/marine-weather/backend && /var/www/marine-weather/venv/bin/python scraper_cron.py >> /var/log/marine-weather-cron.log 2>&1

# ログローテーション - 毎日午前2時
0 2 * * * find /var/log -name "marine-weather-*.log" -mtime +7 -delete